# arch/any/msm-firmware-loader/pkg/msm-firmware-loader/usr/lib/systemd/system/msm-firmware-loader-unpack.service
[Unit]
Description=Unpack Compressed Firmware Blobs
Before=systemd-udev-settle.service
After=msm-firmware-loader.service

[Service]
ExecStart=/usr/bin/msm-firmware-loader-unpack

[Install]
WantedBy=sysinit.target
