# arch/any/msm-firmware-loader/pkg/msm-firmware-loader/usr/lib/systemd/system/msm-firmware-loader.service
[Unit]
Description=Load firmware blobs on qcom devices
Before=systemd-udevd.service
After=local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/msm-firmware-loader

[Install]
WantedBy=sysinit.target
